# config-base/monerod-heartbeat.service
[Unit]
Description=Monerod RPC heartbeat check

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/monerod-heartbeat.sh
